# deploy/personal-budget-pdf.service.example
[Unit]
Description=Personal Budget PDF Streamlit app
After=network.target

[Service]
Type=simple
User=personal-budget
WorkingDirectory=/opt/personal-budget-pdf
EnvironmentFile=/etc/personal-budget-pdf.env
ExecStart=/opt/personal-budget-pdf/.venv/bin/streamlit run app.py --server.address 127.0.0.1 --server.port 8501
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
